# frontend/scripts/com.dprag.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dprag.tunnel</string>

    <!-- 用 autossh 维持隧道；断线自动重连。-M 0 关闭 autossh 自身监控端口，改用 ServerAlive 检测。 -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/autossh</string>
        <string>-M</string>
        <string>0</string>
        <string>-N</string>
        <string>-p</string>
        <string>12638</string>
        <string>-i</string>
        <string>/Users/dp/.ssh/id_ed25519</string>
        <string>-o</string>
        <string>UseKeychain=yes</string>
        <string>-o</string>
        <string>AddKeysToAgent=yes</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=accept-new</string>
        <string>-L</string>
        <string>8000:127.0.0.1:8000</string>
        <string>-L</string>
        <string>8001:127.0.0.1:8001</string>
        <string>-L</string>
        <string>8002:127.0.0.1:8002</string>
        <string>-L</string>
        <string>19530:127.0.0.1:19530</string>
        <string>-L</string>
        <string>3000:127.0.0.1:3000</string>
        <string>root@221.194.152.152</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AUTOSSH_GATETIME</key>
        <string>0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/dprag-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dprag-tunnel.log</string>
</dict>
</plist>
