# scripts/com.dprag.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dprag.frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dp/Desktop/工作文件/DP_rag_skill/scripts/run-frontend.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dp/Desktop/工作文件/DP_rag_skill/frontend</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/dprag-frontend.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dprag-frontend.log</string>
</dict>
</plist>
